# utils/encarne.timer
[Unit]
Description=Checks for new movies to reencode every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
